# supervisord.conf
[supervisord]
nodaemon=true

[program:hubot]
command=/Hubot-DevHub/bin/hubot --adapter devhub --name Hubot
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true

[program:redis]
command=/usr/bin/redis-server
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

